# root/etc/supervisord.d/clash.conf
[program:clash]
command=clash_service
autorestart=true
priority=2
startsecs=5
startretries=3
stopsignal=TERM
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0